# scripts/cluster-kiosk.service
[Unit]
Description=Car Cluster Kiosk Application
After=cluster-update.service network.target
Wants=cluster-update.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/cluster

# X11 표시 환경 설정
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=QT_QPA_PLATFORM=xcb
Environment=QT_QPA_PLATFORMTHEME=
Environment=XDG_RUNTIME_DIR=/run/user/1000

# X 서버가 준비될 때까지 대기 후 실행
ExecStartPre=/bin/bash -c "until xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done"
ExecStart=/opt/cluster/cluster

Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
